# dockerfiles/files/etc/profile.d/common.sh
alias ll='ls -la --color=auto'
alias grep='grep --color'
alias df='df -h'
alias ps='ps --sort=start_time'
